# etc/cron.d/fcch-reboot-nightly
# 04:00 AM nightly
0 4 * * * root /usr/sbin/reboot
